# 99-3isec-tor-hook.rules
SUBSYSTEMS=="xen", KERNEL=="eth*", ACTION=="add", RUN+="/bin/systemctl --no-block 3isec-tor.service"
